# .config/systemd/user/rclone-sync.timer
[Unit]
Description=RClone backup job
Requires=rclone-sync.service

[Timer]
OnBootSec=10m
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target

